# systemd/alt-news.service
[Unit]
Description=ALT 뉴스 수집
After=network.target postgresql.service redis.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
ExecStart=__VENV_PYTHON__ __PROJECT_DIR__/manage.py collect_news
EnvironmentFile=__PROJECT_DIR__/.env
