# systemd/alt-dart.timer
[Unit]
Description=ALT 다트 공시 수집 타이머 (10분)

[Timer]
OnCalendar=*:0/10
Persistent=true

[Install]
WantedBy=timers.target
