# systemd/alt-dart.service
[Unit]
Description=ALT 다트 공시 수집
After=network.target postgresql.service redis.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
ExecStart=__VENV_PYTHON__ __PROJECT_DIR__/manage.py collect_dart
EnvironmentFile=__PROJECT_DIR__/.env
